# proxy/scripts/systemd/mirror-models.timer
[Unit]
Description=Run model mirror periodically

[Timer]
OnBootSec=5min
OnUnitActiveSec=6h
Persistent=true

[Install]
WantedBy=timers.target
